# KI-Evaluation/Szenario A/work/scripts/homelab-backup.service
[Unit]
Description=Run local homelab backups
After=network-online.target pve-cluster.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/homelab-backup.sh
