# KI-Evaluation/Szenario A/work/services/common/homelab-docker-firewall.service
[Unit]
Description=Homelab Docker access restrictions
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/homelab-docker-firewall.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
