# KI-Evaluaion/Szenario A/BTC/homelab-btc-firewall.service
[Unit]
Description=Apply btc-node firewall restrictions
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/homelab-btc-firewall.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
